# s6-rc.d/socat/run
#!/command/execlineb -P

# Run socat to expose joplin clipper api bound to localhost as s6 service
# Forward external port 80 to Joplin server on 127.0.0.1:41184
# We need socat to allow external (0.0.0.0) access to hardcoded localhost port binding:
# https://github.com/laurent22/joplin/blob/d22abe69b649f5909e85a9b72400978980f1f396/packages/lib/ClipperServer.ts#L231

pipeline -w { sed -u "s/^/socat: /" }
foreground { echo "forwarding TCP port 0.0.0.0:80 => 127.0.0.1:41184" }
fdmove -c 2 1
socat -d -ls TCP-LISTEN:80,fork TCP:127.0.0.1:41184
