# s6-rc.d/joplin-config/up
foreground {
    pipeline -w { sed -u "s/^/joplin version: /" }
    fdmove -c 2 1
    joplin version
}
# Merge JSON files in order of precedence and feed result to "joplin config --import"
# (see https://stackoverflow.com/a/71416016/5905029 for jq expression)
multisubstitute {
    importas defaults JOPLIN_CONFIG_DEFAULTS_JSON
    importas config JOPLIN_CONFIG_JSON
    importas required JOPLIN_CONFIG_REQUIRED_JSON
}
pipeline { jq "reduce inputs as \$i (.; . * \$i)" \${defaults} \${config} \${required} }
pipeline {
    fdmove -c 2 1
    joplin config --import
}
sed -u "s/^/joplin config: /"
