# s6-rc.d/joplin-config/up
foreground {
    pipeline -w { sed -u "s/^/joplin version: /" }
    fdmove -c 2 1
    joplin version
}
# Apply the JSON joplin config files in order
pipeline -w { sed -u "s/^/joplin config: /" }
forx -E JSON_FILE {
    JOPLIN_CONFIG_DEFAULTS_JSON
    JOPLIN_CONFIG_JSON
    JOPLIN_CONFIG_REQUIRED_JSON
}
importas json_file $\{JSON_FILE}
if { test -e $\{json_file} }
foreground { echo importing $\{json_file} }
fdmove -c 2 1
joplin config --import-file $\{json_file}
